# root/etc/services.d/timer-for-google-assistant/run
#!/usr/bin/with-contenv bash
npm run --prefix /app start:prod
